# home/.bash/profile.d/prompt.sh
source $HOME/.dotsplat/repos/dotfiles/home/prompt/bash-svn-prompt.sh
source $HOME/.dotsplat/repos/dotfiles/home/prompt/bash-git-prompt.sh
