# home/.bash/profile.d/prompt.sh
PS1="[\t][\u@\h:\w]\$ "
export PS1

source $HOME/.dotsplat/repos/dotfiles/home/.bash/prompt/bash-git-prompt/gitprompt.sh
source $HOME/.dotsplat/repos/dotfiles/home/.bash/prompt/bash-svn-prompt.sh

